# config/launchd/com.tetra.local-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.tetra.local-proxy</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/angwei/Repos/tetra/.venv/bin/python</string>
      <string>/Users/angwei/Repos/tetra/scripts/local_host_proxy.py</string>
      <string>--port</string>
      <string>8080</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/angwei/Repos/tetra</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/angwei/Library/Logs/tetra.local-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/angwei/Library/Logs/tetra.local-proxy.log</string>
  </dict>
</plist>
