# config/launchd/com.tetra.benchmark-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tetra.benchmark-pipeline</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/angwei/Repos/tetra/bin/run_benchmark_pipeline.sh</string>
    </array>
    
    <!-- Run daily at 8:00 PM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/angwei/Repos/tetra</string>
    
    <key>StandardOutPath</key>
    <string>/tmp/tetra-benchmark-pipeline.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/tetra-benchmark-pipeline.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/angwei/.pyenv/shims:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/angwei/Repos/tetra</string>
    </dict>
    
    <!-- Don't run if on battery -->
    <key>RequiresHardpower</key>
    <true/>
    
    <!-- Keep alive during run -->
    <key>KeepAlive</key>
    <false/>
    
    <!-- Nice value for lower priority -->
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>